# ProjectDevOps/roles/supervisor/templates/flask_project.conf
[program:flaskProject]
command = gunicorn hello:app -b localhost:8000
directory = /home/{{ user }}/flaskProject
user = vagrant
